# config/services/api/crontab
* * * * * root /var/www/html/cron.sh >> /var/log/cron.log 2>&1